# seed1.conf
start on runlevel [06]
task
exec /home/seed/src/git/salesapp/mainApp/run_seed.sh
