# mainApp/scripts/seed2.conf
start on startup
task
exec /root/src/git/salesapp/mainApp/run_seed.sh